# deploy/systemd/manga-library-app.service
[Unit]
Description=Manga Library App (node web service)
After=network.target

[Service]
Type=simple
User=svc_manga
Group=svc_manga
WorkingDirectory=/srv/manga-library-app
EnvironmentFile=/etc/manga-library-app/.env
Environment=NODE_ENV=production
Environment=NODE_OPTIONS=--max-old-space-size=512
ExecStart=/usr/bin/env node /srv/manga-library-app/server.js
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

# Notes:
# - Copy this file to /etc/systemd/system/manga-library-app.service and run:
#   systemctl daemon-reload && systemctl enable --now manga-library-app
# - The service expects an EnvironmentFile at /etc/manga-library-app/.env